# linux/udev/99-spacemouse.rules
# CAD Mouse MK2 — VID/PID spoofed as SpaceMouse Compact (256f:c635)
# Restart spacenavd when device is plugged in, then restart user services.
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="256f", ATTRS{idProduct}=="c635", \
    RUN+="/usr/bin/systemctl restart spacenavd", \
    RUN+="/usr/local/bin/restart-spacemouse-services"

# Grant the 'input' group read access to the hidraw interface (for button mapping)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="256f", ATTRS{idProduct}=="c635", \
    MODE="0660", GROUP="input"
